# ansible/roles/fswatch/templates/fswatchers.target.j2
[Unit]
Description=FSwatchers
Requires={{ fswatch_services | join(' ') }}

[Install]
WantedBy=multi-user.target
