# ansible/playbooks/templates/fswatch.service.j2
[Unit]
Description="watch changes in FS at %I"
PartOf=fswatchers.target

[Service]
Restart=always
ExecStart={{ fswatch_caller_script_abspath }} %I
