# ansible/roles/sneakpeek/templates/tcptracer/tcptracer.target.j2
[Unit]
Description="All tcptracer services"
Requires={{ tcptracer_services | join(' ') }}

[Install]
WantedBy=multi-user.target
