# ansible/roles/execsnoop/templates/execsnoops.target.j2
[Unit]
Description="All execsnoop services"
Requires={{ execsnoop_services | join(' ') }}

[Install]
WantedBy=multi-user.target
